# systemd/x570-temp.timer.j2
[Unit]
Description=Regularly read out AMD Ryzen on X570 temperature

[Timer]
OnBootSec=1min
OnUnitActiveSec={{ x570_temp_interval }}

[Install]
WantedBy=timers.target
